# share/systemd/user/home-nor-var.automount
[Unit]
Description=Automount home var overlay

[Automount]
Where=/home/nor/var

[Install]
WantedBy=default.target
